# vps-backup/src/backup.cron
# 每天凌晨2点执行完整备份
0 2 * * * root /usr/local/bin/vps-backup --backup full > /var/log/vps_backup/backup.log 2>&1

# 每6小时执行一次增量备份
0 */6 * * * root /usr/local/bin/vps-backup --backup incremental > /var/log/vps_backup/incremental.log 2>&1
